# ansible/roles/db/templates/postgres-exporter.service.j2
[Unit]
Description=Postgres Exporter
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/postgres_exporter
Restart=always

[Install]
WantedBy=multi-user.target

